# Host/00-backend.rules
ATTRS{devpath}=="1.4", ENV{DEVTYPE}=="disk", SYMLINK+="USBOutputDisk", ATTR{authorized}:="1", RUN+="/bin/mount /dev/%k /var/lib/docker/volumes/USBOutputDevice/_data/USBOutputDisk"
ATTRS{devpath}=="1.4", ENV{DEVTYPE}=="partition", SYMLINK+="USBOutputPart", ATTR{authorized}:="1", RUN+="/bin/mount /dev/%k /var/lib/docker/volumes/USBOutputDevice/_data/USBOutputPart"
ATTRS{devpath}=="1.4", KERNEL=="hid*", SYMLINK+="USBBadOutput", RUN+="/bin/bash -c 'echo 0 > /sys/bus/usb/devices/1-1.4/authorized'", RUN+="/usr/bin/touch /var/lib/docker/volumes/USBOutputDevice/_data/BadUSBAttempt"