# Host/Mounters/inputdiskmnt@.service
[Unit]
Description=Mount USB Drive on %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/share/Asterix/Mounters/inputdiskmnt.sh add %i
ExecStop=/usr/share/Asterix/Mounters/inputdiskmnt.sh remove %i