# Host/Mounters/inputpartmnt@.service
[Unit]
Description=Mount USB Drive on %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/share/Asterix/Mounters/inputpartmnt.sh add %i
ExecStop=/usr/share/Asterix/Mounters/inputpartmnt.sh remove %i