# Host/Mounters/outppartmnt@.service
[Unit]
Description=Mount USB Drive on %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/share/Asterix/mounters/outppartmnt.sh add %i
ExecStop=/usr/share/Asterix/mounters/outppartmnt.sh remove %i