# Host/Mounters/outpdiskmnt@.service
[Unit]
Description=Mount USB Drive on %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/share/Asterix/mounters/outpdiskmnt.sh add %i
ExecStop=/usr/share/Asterix/mounters/outpdiskmnt.sh remove %i